# files/etc/supervisor/conf.d/phpfpm.conf
[program:phpfpm]
command=/usr/sbin/php5-fpm -R --nodaemonize
stopsignal=QUIT
redirect_stderr=true
autorestart=true
